# systemd/bme280-sensor-to-hub.service
[Unit]
Description=BME280 sensor to raspi-esp32-status-panel hub
After=network.target

[Service]
Type=simple

# --- 以下の3行は実際にデプロイしたパスに書き換えてください ---
# 例: /home/pi/bme280_sensor_to_hub
WorkingDirectory=/path/to/bme280_sensor_to_hub
EnvironmentFile=/path/to/bme280_sensor_to_hub/.env
ExecStart=/path/to/bme280_sensor_to_hub/.venv/bin/python -m bme280_sensor_to_hub

# --- サービスを実行するユーザー名に書き換えてください ---
# 例: pi, satomi など (whoami コマンドで確認できます)
User=YOUR_USER

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
